# powermgnt/90-powermgnt-dom0.rules
# handle change of power source ac/bat
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/usr/local/bin/powermgnt-dom0"

